# rootfs/etc/services.d/kube2nginx/run
#!/usr/bin/execlineb -P

with-contenv
s6-envuidgid nginx
multisubstitute
{
  importas -u uid UID
  importas -u gid GID
}
kube2nginx --nginx-src=/etc/nginx/nginx.conf.tmpl --nginx-dst-uid=${uid} --nginx-dst-gid=${gid}
