# rootfs/etc/services.d/kube2nginx/run
#!/usr/bin/execlineb -P

with-contenv
s6-setuidgid nginx
kube2nginx --nginx-template-path=/etc/nginx/nginx.conf.tmpl
